# optimus-manager/optimus-manager.service
[Unit]
Description=Use dedicated graphics card

[Service]
ExecStart=optimus-manager --switch nvidia --no-confirm

[Install]
WantedBy=default.target
